# deployment/config-files/s6-rc.d/oneshot-php-fpm-log-prepare/up
#!/command/execlineb -P

if { mkdir -p /app/var/log/php-fpm }
