# deployment/config-files/s6-rc.d/longrun-supercronic/run
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

/usr/local/bin/supercronic /etc/cron.d/degree_crontab
