# deployment/config-files/s6-rc.d/longrun-nginx/run
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

# stderr and stdout to log
fdmove -c 2 1

nginx
