# deployment/config-files/s6-rc.d/longrun-nginx/run
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

nginx
