# deployment/config-files/s6-rc.d/oneshot-symfony-assets-install/up
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

php /app/bin/console assets:install
