# deployment/config-files/s6-rc.d/longrun-mercure/run
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

/usr/bin/mercure -a 127.0.0.1:3000 --subscriptions --compress
