# deployment/config-files/s6-rc.d/oneshot-wait-for-db/up
#!/command/execlineb -P

with-contenv

bash /etc/s6-overlay/s6-rc.d/oneshot-wait-for-db/wait_for_db.sh
