# deployment/config-files/s6-rc.d/oneshot-supercronic-log-prepare/up
#!/command/execlineb -P

if { mkdir -p /app/var/log/supercronic }
