# deployment/config-files/s6-rc.d/longrun-video-encoder/run
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

foreground {
    sleep 10
}

php /app/bin/console messenger:consume async -vv
